# units/git-eventc-webhook@.service.in
[Unit]
Description=Git WebHook to eventd gateway

[Service]
ExecStart=@bindir@/git-eventc-webhook --token %i

[Install]
WantedBy=default.target
