# units/git-eventc-webhook.service.in
[Unit]
Description=Git WebHook to eventd gateway
Requires=eventd.socket
After=eventd.socket

[Service]
ExecStart=@bindir@/git-eventc-webhook

[Install]
WantedBy=default.target
